# agent/systemd/ops-engine-agent.timer
[Unit]
Description=Run Ops Engine agent every minute

[Timer]
OnBootSec=2min
OnUnitActiveSec=1min
AccuracySec=10s
Persistent=true

[Install]
WantedBy=timers.target
